# templates/init.upstart.erb
task

exec /usr/local/sbin/ipset_init -c "<%= @cfg %>"

# run before all sysv-init scripts, where the firewall (iptables) is
start on starting rc RUNLEVEL=[2345]
