# templates/systemd.service.erb
[Unit]
Description=define and fill-in ipsets
Documentation=https://github.com/mighq/puppet-ipset
Before=<%= @before_target %>

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/ipset_init -c "<%= @cfg %>"

[Install]
WantedBy=basic.target
