# templates/systemd.service.erb
[Unit]
Description=Define ipsets
Documentation=https://github.com/mighq/puppet-ipset
Before=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/ipset_init -c "<%= @cfg %>"

[Install]
WantedBy=basic.target
